# backend/scripts/activity-reporter.service
[Unit]
Description=Neural Grid Activity Reporter for %i
After=network.target

[Service]
Type=simple
ExecStart=/opt/neural-grid/activity-reporter.sh %i http://100.114.123.105:3101
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
